# deploy/hermes/systemd/hermes-backup.service
# Template — not installed by HERMES-DEP-001/HERMES-INSTALL-001-preparation.
# Review before copying to /etc/systemd/system/.
#
# One-shot backup job, triggered by hermes-backup.timer. Runs a host-level
# tar of the single shared Hermes data volume (/opt/hermes/data) —
# deploy/hermes/scripts/run-backup.sh. No docker exec, no `docker` group
# membership required for this job (unlike hermes-healthcheck.service,
# which does need it — see that unit's comment): /opt/hermes/data is a
# plain host directory the `hermes` user already owns.

[Unit]
Description=Hermes backup job (host-level tar of the shared data volume)
Documentation=file:///opt/hermes/runbooks/BACKUP_RESTORE.md

[Service]
Type=oneshot
User=hermes
Group=hermes
WorkingDirectory=/opt/hermes
ExecStart=/opt/hermes/core/compose/scripts/run-backup.sh

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/hermes/backups /opt/hermes/logs
ReadOnlyPaths=/opt/hermes/data
ProtectHome=yes
