# deploy/hermes/systemd/hermes-healthcheck.service
# Template — not installed by HERMES-DEP-001. Review before copying to /etc/systemd/system/.
#
# One-shot liveness check, triggered by hermes-healthcheck.timer. Only
# observes and logs — it never restarts anything itself; recovery is
# Docker's own `restart: unless-stopped` policy on the ict-trading
# container (see deploy/hermes/core/docker-compose.yml.template — there is
# no custom gateway-supervising systemd unit in this package). See
# deploy/hermes/runbooks/HEALTH_CHECKS.md.
#
# Same docker-group trade-off as hermes-backup.service — see that unit's
# comment.

[Unit]
Description=Hermes health check (docker compose exec hermes doctor/status)
Documentation=file:///opt/hermes/runbooks/HEALTH_CHECKS.md
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=hermes
Group=hermes
WorkingDirectory=/opt/hermes/core/compose
ExecStart=/opt/hermes/core/compose/scripts/run-healthcheck.sh

NoNewPrivileges=yes
ProtectSystem=strict
ReadOnlyPaths=/opt/hermes/core/compose
ReadWritePaths=/opt/hermes/logs
ProtectHome=yes
