# deploy/hermes/systemd/hermes-worker.service
# Template — not installed by HERMES-DEP-001. Review before copying to /etc/systemd/system/.
#
# Supervises the optional Hermes worker Compose service (docs/10_INFRASTRUCTURE.md
# — Services: "optional worker service"). Depends on hermes-gateway.service for
# the shared Compose network/volumes.

[Unit]
Description=Hermes worker (Docker Compose service)
Documentation=file:///opt/hermes/README.md
After=hermes-gateway.service
Requires=hermes-gateway.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/hermes/core/compose
EnvironmentFile=/opt/hermes/core/secrets/.env
ExecStart=/usr/bin/docker compose up -d worker
ExecStop=/usr/bin/docker compose stop worker
Restart=on-failure
RestartSec=10

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/hermes/core/compose /opt/hermes/logs
ProtectHome=yes

[Install]
WantedBy=multi-user.target
